# .config/systemd/user/wallpaper.service
[Unit]
Description=Change wallpapers using wallpaperscript.sh

[Service]
ExecStart="/home/deus/.local/bin/wallpaperscript.sh"
